# systemd/user/stop-and-think.timer
[Unit]
Description=Stop and think every 90 minutes

[Timer]
OnBootSec=60min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
