# systemd/user/stop-and-think.timer
[Unit]
Description=Stop and think every 90 minutes

[Timer]
OnBootSec=90min
OnUnitActiveSec=90min

[Install]
WantedBy=timers.target
